# deploy/promanager.service
[Unit]
Description=ProManager (FastAPI + static UI)
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/promanager/backend
Environment=PYTHONUNBUFFERED=1
# 按实际路径修改 /opt/promanager
ExecStart=/opt/promanager/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
